# ansible/roles/restic/templates/restic@.timer.j2
[Unit]
Description=Run Restic backup on schedule

[Timer]
OnCalendar={{ restic_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
